# scripts/install/conf/rq-all.target
[Unit]
Description=Workers
Requires=rq@1.service rq@2.service rq@3.service rq@4.service

[Install]
WantedBy=multi-user.target